# cuca.service
[Unit]
After=etcd

[Service]
ExecStart=/usr/bin/beacon
Type=simple

